# instagram-gemini-bot/deploy/instagram-bot.service
[Unit]
Description=Instagram Gemini Bot (webhook)
After=network.target

[Service]
Type=simple
# O'zingizga moslang: bot papkasi qayerda joylashgan
WorkingDirectory=/opt/instagram-gemini-bot
# .env fayldan sozlamalarni o'qiydi
EnvironmentFile=/opt/instagram-gemini-bot/.env
# gunicorn bilan 2 ta worker, 5000-portда
ExecStart=/opt/instagram-gemini-bot/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
